# archived/201812/nettool/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D


[program:shadowsocks]
command=/usr/bin/ssserver -c /shadowsocks.json
