# archived/201812/dojo-doc/server.conf
[program:www]
command=/bin/bash -c "cd /dojo;make run_httpserver"


[program:wiki]
command=/bin/bash -c "cd /dojo;make run_wiki"

