# archived/201812/dojo-doc/supervisord.conf
[program:rstwiki]
command=/usr/bin/make run_wiki
[program:dapi]
command=/usr/bin/make run_dapi
[program:httpserver]
command=/usr/bin/make run_httpserver

